# docs/examples/systemd/chv-stord.service
[Unit]
Description=CHV Storage Daemon
After=network.target

[Service]
Type=simple
User=chv
Group=chv
ExecStartPre=/bin/mkdir -p /run/chv/stord
ExecStart=/usr/local/bin/chv-stord /etc/chv/stord.toml
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=5
RuntimeDirectory=chv/stord
StateDirectory=chv
LogsDirectory=chv
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/chv/stord /var/lib/chv/storage

[Install]
WantedBy=multi-user.target
