# docs/examples/systemd/chv-agent.service
[Unit]
Description=CHV Node Agent
After=network.target chv-controlplane.service
Wants=chv-controlplane.service

[Service]
Type=simple
User=chv
Group=chv
UMask=002
ExecStartPre=/bin/mkdir -p /run/chv/agent /run/chv/stord /run/chv/nwd /var/lib/chv/agent
ExecStart=/usr/local/bin/chv-agent /etc/chv/agent.toml
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=5
RuntimeDirectory=chv
StateDirectory=chv
LogsDirectory=chv

# Device access for KVM
DeviceAllow=/dev/kvm rw
SupplementaryGroups=kvm

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/chv /var/lib/chv
ReadOnlyPaths=/etc/chv /usr/local/bin /usr/bin

# Clean up stale sockets if the agent exits uncleanly
ExecStopPost=-/bin/rm -f /run/chv/agent/api.sock

[Install]
WantedBy=multi-user.target
