# docs/examples/systemd/chv-nwd.service
[Unit]
Description=CHV Network Daemon
After=network.target

[Service]
Type=simple
User=chv
Group=chv
ExecStart=/usr/local/bin/chv-nwd /etc/chv/nwd.toml
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=5
RuntimeDirectory=chv/nwd
StateDirectory=chv
LogsDirectory=chv
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/run/chv/nwd /run/netns

[Install]
WantedBy=multi-user.target
