# legacy/go-controlplane/ops/chv-controller.service
[Unit]
Description=CHV Controller - Cloud Hypervisor Virtualization Controller
After=network.target chv-agent.service
Requires=chv-agent.service

[Service]
Type=simple
ExecStart=/usr/local/bin/chv-controller
Restart=on-failure
RestartSec=5
User=root
Environment="CHV_HTTP_ADDR=:8888"
Environment="CHV_DATA_ROOT=/var/lib/chv"
Environment="CHV_AGENT_URL=http://localhost:9090"
Environment="CHV_LOG_LEVEL=info"
Environment="CHV_LOG_DIR=/var/lib/chv/logs"

[Install]
WantedBy=multi-user.target
